# Linux/vmware-tools.service
[Unit]
Description=VMWare Tools Service

[Service]
ExecStart=/home/voidflow/Autostart/vmware.sh

[Install]
WantedBy=multi-user.target
